# apt-daily/apt-update.service
[Unit]
Description=apt-get update Service

[Service]
Type=simple
ExecStart=/usr/bin/apt-get update


